# deploy/plesk/staging/systemd/crigestion-staging-recovery-bundle.service
[Unit]
Description=Create encrypted CriGestion staging recovery bundle
After=network-online.target postgresql.service crigestion-staging-backup.service
Wants=network-online.target
OnFailure=crigestion-staging-recovery-bundle-alert.service

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/crigestion-staging/current
Environment=NODE_ENV=production
EnvironmentFile=/etc/crigestion-staging/recovery-bundle.env
Environment=RECOVERY_BUNDLE_WORK_DIRECTORY=/run/crigestion-staging-recovery
LoadCredential=recovery-bundle.key:/etc/crigestion-staging/recovery-bundle.key
RuntimeDirectory=crigestion-staging-recovery
RuntimeDirectoryMode=0700
ExecStart=/usr/local/sbin/crigestion-staging-recovery-bundle
TimeoutStartSec=30min
UMask=0077

NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
ProtectHome=read-only
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native
ReadWritePaths=/root/crigestion-staging-recovery
InaccessiblePaths=/etc/crigestion-staging/recovery-bundle.key

StandardOutput=journal
StandardError=journal
