# deploy/plesk/staging/systemd/crigestion-staging-recovery-drill.service
[Unit]
Description=Validate a CriGestion staging recovery bundle in isolation
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/crigestion-staging/current
Environment=NODE_ENV=production
Environment=RECOVERY_BUNDLE_WORK_DIRECTORY=/var/lib/crigestion-staging-recovery-drill
LoadCredential=recovery-bundle.key:/etc/crigestion-staging/recovery-bundle.key
StateDirectory=crigestion-staging-recovery-drill
StateDirectoryMode=0700
ExecStartPre=/usr/bin/test -x /usr/bin/python3
ExecStartPre=/usr/bin/test -f /opt/crigestion-staging/current/scripts/extract-recovery-bundle.py
ExecStart=/usr/local/sbin/crigestion-staging-recovery-drill
TimeoutStartSec=60min
UMask=0077

NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
ProtectHome=read-only
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native
ReadWritePaths=/var/lib/crigestion-staging-recovery-drill
InaccessiblePaths=/etc/crigestion-staging/recovery-bundle.key
MemoryMax=1G

StandardOutput=journal
StandardError=journal
