# deploy/plesk/staging/systemd/crigestion-staging-backup.timer
[Unit]
Description=Run CriGestion staging PostgreSQL backup daily

[Timer]
OnCalendar=*-*-* 02:15:00
RandomizedDelaySec=15min
Persistent=true
Unit=crigestion-staging-backup.service

[Install]
WantedBy=timers.target
