# deploy/plesk/staging/systemd/crigestion-staging-app.service
[Unit]
Description=CriGestion staging web application
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=crigestion-staging
Group=crigestion-staging
WorkingDirectory=/opt/crigestion-staging/current
EnvironmentFile=/etc/crigestion-staging/app.env
Environment=NODE_ENV=production
Environment=HOSTNAME=127.0.0.1
Environment=PORT=3101
ExecStartPre=/usr/bin/test -f /opt/crigestion-staging/current/.next/BUILD_ID
ExecStart=/opt/plesk/node/22/bin/node node_modules/next/dist/bin/next start --hostname 127.0.0.1 --port 3101
Restart=on-failure
RestartSec=5s
TimeoutStartSec=60s
TimeoutStopSec=45s
KillSignal=SIGTERM
UMask=0027

NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectControlGroups=true
ProtectHome=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native
ReadWritePaths=/var/cache/crigestion-staging

[Install]
WantedBy=multi-user.target
